# deploy/ci/upstart.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 3 5

setuid respect
setgid respect

exec /home/respect/ci/deploy/ci/service_demon.sh