# deploy/production/upstart.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 3 5

setuid respect
setgid respect

exec /home/respect/production/deploy/production/service_demon.sh